# .docker/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:dobby]
priority=10
directory=/dobby
command=pipenv run python dobby run config.yml
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0